# .docker/aio/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/17/bin/postgres -D %(ENV_PGDATA)s
user=postgres
priority=10
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:heya]
command=/usr/local/bin/heya-aio-start-heya
priority=20
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:heya-worker]
command=/usr/local/bin/heya-aio-start-worker
# Let the API process run migrations and bring Caddy up first; queue recovery
# can then take as long as necessary without delaying ingress readiness.
priority=30
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
